# imageroot/systemd/user/nethvoice-hotel-alarms.timer
[Unit]
Description=Run NethVoice hotel alarm calls every minute

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
